# framework/inventory/awscli/rootfs/etc/profile.d/20-autocomplete.sh
#!/usr/bin/env bash
case $- in
  *i*) [ -z "$PS1" ] && return ;;
  *) return ;;
esac

complete -C "/usr/local/aws-cli/v2/current/bin/aws_completer" aws
