# framework/module/abstract/score/provision/rootfs/etc/profile.d/0-environment.sh
shopt -os allexport && {
  . /etc/environment
  shopt -ou allexport
}


