# framework/module/abstract/score/provision/rootfs/etc/profile.d/0-environment.sh
if [ -z "$PS1" ]; then
  return
fi

shopt -os allexport && {
  . /etc/environment
  shopt -ou allexport
}


