# framework/module/abstract/score/provision/rootfs/etc/profile.d/0-environment.sh
case $- in
  *i*) [ -z "$PS1" ] && return ;;
  *) return ;;
esac

shopt -os allexport && {
  . /etc/environment
  shopt -ou allexport
}


